# nginx/GetData.service
[Unit]
 Description=Service description

[Service]
 ExecStart=/bin/get-data.py

[Install]
 WantedBy=default.target
